# img/units/write-issue.service
[Unit]
Description=Write a useful /etc/issue

[Service]
Type=oneshot
ExecStart=/usr/local/bin/write-issue

[Install]
WantedBy=multi-user.target
